# supervisord.conf
[supervisord]
nodaemon=true
loglevel=critical

[program:cron]
command = cron -f -L 15
autostart=true
autorestart=true

[program:haproxy]
command = /certs.sh && /docker-entrypoint.sh haproxy -f /usr/local/etc/haproxy/haproxy.cfg
autostart=true
autorestart=true
